# daily-news-bot.service
[Unit]
Description=Daily News Discord Bot
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/daily-news-bot/publish
ExecStart=/usr/bin/dotnet /opt/daily-news-bot/publish/DailyNewsBot.dll
Restart=on-failure
RestartSec=15
StartLimitBurst=3
StartLimitIntervalSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-news-bot

[Install]
WantedBy=multi-user.target
